# app/control-api/systemd/stemford-memory-cards-maintenance.service
[Unit]
Description=Stemford memory cards maintenance
After=network-online.target stemford-control-api.service
Wants=network-online.target

[Service]
Type=oneshot
User=stemford
Group=stemford
WorkingDirectory=/opt/stemford/app/control-api
EnvironmentFile=/opt/stemford/run/.env
ExecStart=/opt/stemford/app/control-api/scripts/memory_cards_maintenance.sh
NoNewPrivileges=true

